# services/norm-face.service
[Unit]
Description=N.O.R.M. Face Core
After=network-online.target sound.target
Wants=network-online.target

[Service]
User=norm
WorkingDirectory=/home/norm/norm
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=SDL_RENDER_DRIVER=software
Environment=SDL_AUDIODRIVER=dummy
Environment=NORM_TOUCH_DEVICE=/dev/input/event0
ExecStart=/home/norm/norm/.venv/bin/python /home/norm/norm/app.py --config /home/norm/norm/configs/norm-alpha.json
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
